# treatwell-outreach.service
[Unit]
Description=Treatwell Outreach — Daily Email Send
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/mma/treatwell-outreach
ExecStart=/home/mma/treatwell-outreach/run_pipeline.sh
StandardOutput=append:/home/mma/treatwell-outreach/logs/pipeline.log
StandardError=append:/home/mma/treatwell-outreach/logs/pipeline.log
Environment=PYTHONUNBUFFERED=1
TimeoutStartSec=3600

[Install]
WantedBy=default.target
